# wa_service/voidpanel-wa.service
[Unit]
Description=VoidPanel WhatsApp Web Microservice (Baileys)
Documentation=https://github.com/WhiskeySockets/Baileys
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/panel/wa_service
ExecStart=/usr/local/bin/node /var/www/panel/wa_service/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voidpanel-wa
Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin
# Ensure service only listens on localhost — Django proxies to it
Environment=HOST=127.0.0.1
Environment=PORT=3001

[Install]
WantedBy=multi-user.target
